# package/base-files/files/etc/init.d/z_custom
#!/bin/sh /etc/rc.common

START=99
STOP=90
USE_PROCD=1

service_triggers()
{
	procd_add_reload_trigger "system"
}

start_service()
{
	sleep 20s && /etc/init.d/ttyd restart &
}
