# target/linux/rockchip/armv8/base-files/etc/init.d/set-irq-affinity
#!/bin/sh /etc/rc.common

START=99

get_device_irq() {
	local device="$1"
	local line
	local seconds="0"

	# wait up to 10 seconds for the irq/device to appear
	while [ "${seconds}" -le 10 ]; do
		line=$(grep -E -m 1 "${device}\$" /proc/interrupts) && break
		seconds="$(( seconds + 2 ))"
		sleep 2
	done
	echo ${line} | sed 's/:.*//'
}

set_interface_core() {
	local core_mask="$1"
	local device="$2"

	local irq=$(get_device_irq "$device")

	echo "${core_mask}" > /proc/irq/${irq}/smp_affinity
}

start() {
	uci set network.globals.packet_steering='0'
	uci commit network
	
	echo "32768" > "/proc/sys/net/core/rps_sock_flow_entries"
	set_interface_core 8 "eth0"
	set_interface_core 4 "eth1"
	set_interface_core 8 "xhci-hcd:usb7"
	set_interface_core 8 "xhci-hcd:usb9"
	echo "3" > "/sys/class/net/eth0/queues/rx-0/rps_cpus"
	echo "3" > "/sys/class/net/eth1/queues/rx-0/rps_cpus"
	echo "4096" > "/sys/class/net/eth0/queues/rx-0/rps_flow_cnt"
	echo "4096" > "/sys/class/net/eth1/queues/rx-0/rps_flow_cnt"
	echo "3" > "/sys/class/net/eth0/queues/tx-0/xps_cpus"
	echo "3" > "/sys/class/net/eth1/queues/tx-0/xps_cpus"

	sysctl -w net.core.netdev_max_backlog=2000
}
